# etc/udev/rules.d/00-usb-keyboards.rules
ACTION=="add", ATTRS{idVendor}=="0c45", ATTRS{idProduct}=="086a", RUN+="/usr/    local/bin/init-keyboard", OWNER="al"
